# templates/org.jenkins-ci.slave.jnlp.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-/Apple/DTD PLIST 1.0/EN" "http:/www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.jenkins-ci.slave.jnlp</string>
    <key>ProgramArguments</key>
    <array>
        <string><%= slave_home %>/run-jenkins-slave</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>UserName</key>
    <string><%= slave_user %></string>
    <key>WorkingDirectory</key>
    <string><%= slave_home %></string>
    <key>SessionCreate</key>
    <true/>
    <key>StandardInPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string><%= log_path %>/org.jenkins-ci.slave.jnlp.err</string>
    <key>StandardOutPath</key>
    <string><%= log_path %>/org.jenkins-ci.slave.jnlp.out</string>
</dict>
</plist>
